# .cron/crontab
BACKUP=/home/sorucoder/.backup
BACKUP_HOST=sorucoder@sorucoder-server
0 2 * * * /home/sorucoder/.cron/backup.sh > /dev/null 2>&1
