# .cron/crontab
0 2 * * * /home/sorucoder/.cron/backup.sh > /dev/null 2>&1
0 4 * 1 * /home/sorucoder/.cron/update.sh > /dev/null 2>&1
