# .cron/crontab
0 2 * * * chronic /home/sorucoder/.cron/backup.sh
0 4 * * 1 chronic /home/sorucoder/.cron/update.sh
